# ansible/roles/restic_backup/templates/s1-backup.service.j2
# {{ ansible_managed }}
[Unit]
Description=Weekly restic backup of s1 data to Cloudflare R2
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
EnvironmentFile={{ restic_env_file }}
ExecStart={{ restic_script_path }}
